# modules/bobwebui/bobwebui.service
[Unit]
Description=Control systemd services through Web or REST API
Documentation=https://github.com/ogarcia/sysdweb
After=network.target
Requires=dbus.socket

[Service]
WorkingDirectory=${BOB_ROOT}/modules/bobwebui
ExecStart=${BOB_ROOT}/modules/bobwebui/.venv/bin/python3 -m src.main -c src/bobwebui.conf
Restart=on-failure

[Install]
WantedBy=multi-user.target
